# ICantBelieveItsNotValetudo/root/etc/services.d/ICantBelieveItsNotValetudo/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info 'Start ICantBelieveItsNotValetudo - kubiq Renderer'
cd /app
npm start
